# etc/systemd/system/mysql.service.d/override-mysql-leader.service
[Service]
EnvironmentFile=/etc/mysql-leader.env
ExecStartPre=/usr/local/bin/mysql-set-readonly -service-pre
